# vm_scripts/tc.service
[Unit]
After=network.target

[Service]
ExecStart=/usr/local/bin/tc-latency.sh

[Install]
WantedBy=default.target
